# deploy/systemd/al-apt-cache-clean.service
[Unit]
Description=AL apt cache cleanup

[Service]
Type=oneshot
ExecStart=__APP_DIR__/scripts/clean-apt-cache.sh
