# deploy/systemd/al-backend.service
[Unit]
Description=AL Backend API
After=network-online.target mongod.service
Wants=network-online.target

[Service]
Type=simple
User=__APP_USER__
Group=__APP_USER__
WorkingDirectory=__APP_DIR__/apps/backend
EnvironmentFile=/etc/al/backend.env
Environment=HOME=__APP_ROOT__
Environment=UV_PYTHON_INSTALL_DIR=__APP_ROOT__/python
ExecStart=__APP_DIR__/apps/backend/.venv/bin/uvicorn al_backend.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
